# etc/xdg/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Clock and battery
Exec=conky -d -c /etc/xdg/conky/conkyrc
Terminal=false
