# etc/xdg/autostart/feh.desktop
[Desktop Entry]
Type=Application
Name=Wallpaper
Exec=feh --bg-scale /etc/nixos/noppes/wallpaper.jpg
Terminal=false
