# etc/xdg/autostart/terminal.desktop
[Desktop Entry]
Type=Application
Name=Terminal
Exec=xfce4-terminal
Terminal=false
