# etc/xdg/autostart/volumeicon.desktop
[Desktop Entry]
Type=Application
Name=Volume Icon
Exec=volumeicon
Terminal=false
